# SupportingScripts/etc/systemd/system/fcgimono.socket
[Unit]
Description=fcgi mono socket

[Socket]
ListenStream=/tmp/fcgimono.sock
Accept=false


